# apps/deeplearning/LLM/5090/run_gemma4_ollama/server/gemma4-ollama-8082.service
[Unit]
Description=Gemma4 Ollama status and proxy service on port 8082
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/Users/tinyos/devel_opment/BerePi/apps/deeplearning/LLM/5090/run_gemma4_ollama/server
Environment=OLLAMA_MODEL=gemma4:31b
Environment=OLLAMA_CONTEXT_LENGTH=8192
Environment=OLLAMA_BASE_URL=http://127.0.0.1:11434
Environment=OLLAMA_HOST=127.0.0.1:11434
Environment=GEMMA4_SERVER_HOST=0.0.0.0
Environment=GEMMA4_SERVER_PORT=8082
Environment=AUTO_PULL=1
ExecStart=/Users/tinyos/devel_opment/BerePi/apps/deeplearning/LLM/5090/run_gemma4_ollama/server/run_service.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
